# infra/euroleague-proxy/cloudflared.service
[Unit]
Description=Cloudflare Quick Tunnel for Euroleague proxy
After=network.target tinyproxy.service
OnFailure=notify-failure@%n.service

[Service]
ExecStart=/opt/euroleague-proxy/start-tunnel.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
