# infra/euroleague-proxy/notify-failure@.service
[Unit]
Description=Telegram failure notification for %i

[Service]
Type=oneshot
ExecStart=/usr/bin/gh workflow run notify_proxy_failure.yaml \
  --repo Maccabipedia/maccabipedia_mediawiki_bot \
  --field service=%i
